# python/unfollowers/crontab.sh.crontab.sample
# change username to your github username
0 0 * * * ./fetch-and-historize.sh username
